# hub/deploy/systemd/termipod-host.service
# Termipod host-agent — systemd unit
#
# Pairs with termipod-hub.service. Same install pattern, different role:
# this runs the daemon that launches backend CLIs (claude, codex) inside
# tmux panes on this machine.
#
# Install:
#   # 1. Build (on any Linux box with Go 1.23+; cross-compile if target differs).
#   cd hub && go build -o /tmp/host-agent ./cmd/host-agent
#
#   # 2. System user with a real $HOME — backend CLIs expect to read
#   #    ~/.claude/, ~/.codex/, ~/.gitconfig, ~/.ssh/ etc.
#   sudo useradd --system --create-home --home /var/lib/termipod-host \
#                --shell /bin/bash termipod-host
#   sudo install -o root -g root       -m 0755 /tmp/host-agent /usr/local/bin/host-agent
#
#   # 3. Token from the hub box:  hub-server tokens issue -kind host …
#   sudo install -d -o root -g termipod-host -m 0750 /etc/termipod-host
#   printf 'HUB_URL=https://hub.example.com\nHUB_TEAM=default\nHUB_TOKEN=paste-here\n' \
#       | sudo install -o root -g termipod-host -m 0640 /dev/stdin /etc/termipod-host/env
#
#   # 4. Activate.
#   sudo install -m 0644 termipod-host.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now termipod-host
#   sudo journalctl -fu termipod-host
#
# Hardening note: this unit is deliberately looser than termipod-hub.service
# because the daemon shells out to Node.js-based CLIs (claude-code) that need
# JIT — MemoryDenyWriteExecute and ProtectHome would break them. Tighten only
# if you know your backend tolerates it.

[Unit]
Description=Termipod host-agent (launches backend CLIs in tmux)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=termipod-host
Group=termipod-host
EnvironmentFile=/etc/termipod-host/env
ExecStart=/usr/local/bin/host-agent run \
          --hub ${HUB_URL} \
          --team ${HUB_TEAM} \
          --token ${HUB_TOKEN} \
          --name %H \
          --launcher tmux \
          --tmux-session hub-agents
Restart=on-failure
RestartSec=3s

NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
